# system/mediaauto-panel.service
[Unit]
Description=MediaAuto Panel (FastAPI)
After=network.target

[Service]
Type=simple
User=mediaauto
WorkingDirectory=/opt/mediaauto
ExecStart=/opt/mediaauto/venv/bin/uvicorn src.panel:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment=PYTHONPATH=/opt/mediaauto

[Install]
WantedBy=multi-user.target
